# packaging/systemd/lumacore-daemon.service
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=LumaCore RGB backend daemon
After=local-fs.target

[Service]
Type=simple
Group=lumacore
UMask=0007
ExecStart=/usr/bin/lumacore-daemon --socket /run/lumacore/lumacore.sock
RuntimeDirectory=lumacore
RuntimeDirectoryMode=0775
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
